# projects/miner/SOURCES/miner.cron
# archive the miner jobs data

* * 1 * * root miner-manage db archive